# files/systemd/mfd-daily-reboot.timer
[Unit]
Description=Run daily MFD kiosk reboot

[Timer]
OnCalendar=*-*-* ${REBOOT_TIME}:00
Persistent=true

[Install]
WantedBy=timers.target
