# files/systemd/mfd-browser-healthcheck.service
[Unit]
Description=Restart the kiosk session if the browser is not running

[Service]
Type=oneshot
ExecStart=/usr/bin/bash -c 'pgrep -u ${KIOSK_USER} -f ${BROWSER_BIN} >/dev/null || systemctl restart getty@tty1.service'
