# files/systemd/mfd-daily-reboot.service
[Unit]
Description=Daily reboot for MFD kiosk

[Service]
Type=oneshot
ExecStart=/usr/sbin/reboot
